# dbus-testbed/container1-usbcreator/usb-creator-dbus.service
[Unit]
Description=Simulated USBCreator D-Bus Service (CVE-2013-1060)
After=dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/usb-creator-service.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
