# dbus-testbed/container3-p2p-signal/powermanager1-signal.service
[Unit]
Description=PowerManager1 Signal Listener (PoC 5 Target)
After=dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/signal-listener-service.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
